# ansible/roles/prometheus/templates/prometheus.service.j2
; DO NOT EDIT ON SERVER
; THIS FILE IS MANAGED BY ANSIBLE

[Unit]
Description=Prometheus service
After=network.target
Requires=docker.service

[Service]
Type=simple
ExecStart=/usr/bin/docker run \
  --network=host \
  --user '{{ prometheus.uid }}:{{ prometheus.gid }}' \
  --mount type=bind,source=/etc/prometheus,target=/etc/prometheus \
  --mount type=bind,source=/var/prometheus,target=/var/prometheus \
  prom/prometheus:v{{ prometheus.version }} \
  --config.file=/etc/prometheus/prometheus.yaml \
  --storage.tsdb.path=/var/prometheus
Restart=always

[Install]
WantedBy=multi-user.target
