# ansible/roles/node_exporter/templates/node-exporter.service
; DO NOT EDIT ON SERVER
; THIS FILE IS MANAGED BY ANSIBLE

[Unit]
Description=node-exporter service
After=network.target
Requires=docker.service

[Service]
Type=simple
ExecStart=/usr/local/bin/node_exporter
Restart=always

[Install]
WantedBy=multi-user.target
